# etc/cron.d/minecraft
0 4 * * * debian /home/debian/bin/minecraft_backup.sh

